# playwright-exporter.service
[Unit]
Description=Playwright Prometheus Exporter
Documentation=https://github.com/maravexa/playwright-exporter
After=network.target

[Service]
Type=simple
User=playwright-exporter
Group=playwright-exporter
ExecStart=/usr/local/bin/playwright-exporter -config /etc/playwright-exporter/config.yaml
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/playwright-browsers
Restart=on-failure
RestartSec=5s
MemoryMax=2G
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/tmp/playwright-exporter /opt/playwright-browsers /opt/playwright-tests
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
